# ops/systemd/gauntlet-controller.timer
[Unit]
Description=Run Gauntlet midnight close and two frozen-evidence retries

[Timer]
OnCalendar=*-*-* 00:00:00 America/Chicago
OnCalendar=*-*-* 00:15:00 America/Chicago
OnCalendar=*-*-* 00:30:00 America/Chicago
AccuracySec=1s
RandomizedDelaySec=0
Persistent=false
Unit=gauntlet-controller.service

[Install]
WantedBy=timers.target
